# grafanazbi/tree/etc/cron.d/import_crons
# Import crons from /localsharec
*/5 * * * * root cp -r /localshared/etc/cron.d/* /etc/cron.d/ && service cron reload
