# grafanazbi/tree/etc/services/hive-lookup/run
#!/bin/bash

cd /localshared/hive-lookup

exec python app.py &>> app.log
